# desktop-overlay/etc/cron.daily/clearcache
#!/bin/bash
sync && sysctl -w vm.drop_caches=3