# desktop-overlay/etc/cron.daily/prelink
#!/bin/bash
[[ -x /usr/bin/prelink ]] && /usr/bin/prelink -amfR &>/dev/null