# jwm-overlay/etc/crontab
0 091420 * * * ~/.update/update-notifier