# jwm-overlay/etc/crontab
* * /14 * * * root /usr/bin/clearcache
* */3 * * * update-notifier
