# jwm-overlay/etc/crontab
#update-notifier
* */3 * * * /usr/bin/update-notifier


