# jwm-overlay/etc/cron.daily/prelink
#!/bin/bash
[[ -x /usr/bin/prelink ]] && 
/usr/bin/prelink -amR &>/dev/null