# jwm-overlay/etc/cron.daily/prelink.cron
#!/bin/bash
[[ -x /usr/sbin/prelink ]] && /usr/sbin/prelink -amR &>/dev/null